# deploy/com.openai.ai-cio.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.openai.ai-cio</string>
  <key>ProgramArguments</key><array>
    <string>REPLACE_WITH_VENV_BIN/cio</string>
    <string>daily-review</string><string>--account-label</string><string>Agentic</string>
  </array>
  <key>WorkingDirectory</key><string>REPLACE_WITH_REPOSITORY_PATH</string>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>REPLACE_WITH_REPOSITORY_PATH/outputs/scheduler.log</string>
  <key>StandardErrorPath</key><string>REPLACE_WITH_REPOSITORY_PATH/outputs/scheduler-error.log</string>
</dict></plist>
